# etc/profile.d/defaults.sh
export EDITOR="$(if [[ -n $DISPLAY ]]; then echo 'kate'; else echo 'joe'; fi)"